# autobuild-tridactyl.service
[Unit]
Description=autobuild-tridactyl

[Service]
ExecStart=%h/projects/buildbot/autobuild.sh
