# autobuild-tridactyl.service
[Unit]
Description=autobuild-tridactyl

[Service]
ExecStart=%h/projects/tridactyl/autobuild-tridactyl.sh
